# labelle-web.service
[Unit]
Description=Labelle Web
After=network.target

[Service]
WorkingDirectory=/opt/labelle-web
ExecStart=/opt/labelle-web/.venv/bin/python server/app.py
Restart=on-failure
RestartSec=5
EnvironmentFile=-/opt/labelle-web/.env

[Install]
WantedBy=multi-user.target
